# launchd/_template.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>__LABEL__</string>
  <key>ProgramArguments</key>
  <array>
    <string>__ALFRED_BIN__/agent-launch</string>
    <string>__SCRIPT__</string>
  </array>
__SCHEDULE_BLOCK__
  <key>StandardOutPath</key>
  <string>/tmp/__LOG_STEM__.stdout</string>
  <key>StandardErrorPath</key>
  <string>/tmp/__LOG_STEM__.stderr</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__PATH__</string>
__JAVA_BLOCK__
    <key>HOME</key>
    <string>__HOME__</string>
    <key>ALFRED_HOME</key>
    <string>__ALFRED_HOME__</string>
    <key>ALFREDRC</key>
    <string>__ALFREDRC__</string>
    <key>WORKSPACE_ROOT</key>
    <string>__WORKSPACE_ROOT__</string>
    <key>AGENT_CODENAME</key>
    <string>__AGENT_SHORT__</string>
    <key>LAUNCHD_LABEL</key>
    <string>__LABEL__</string>
__GH_ORG_BLOCK__
__ROLE_BLOCK__
  </dict>
  <key>WorkingDirectory</key>
  <string>__HOME__</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
